# systemd/backup-script.timer
[Unit]

[Timer]
OnCalendar=weekly
RandomizedDelaySec=1d
Unit=backup-script.service

[Install]
WantedBy=timers.target
