# systemd/backup-script.service
[Unit]

[Service]
ExecStart=/opt/root/Backup-Script/run-backup.sh
Type=simple
User=root
